# packaging/systemd/komp-timetracker.service
[Unit]
Description=Komp TimeTracker - Parental Control Service
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 -m komp_timetracker.service
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1

# Configuration paths
Environment=KOMP_CONFIG=/etc/komp-timetracker/config.yaml
Environment=KOMP_DB=/var/lib/komp-timetracker/usage.db

# Required permissions
Capabilities=CAP_SYS_ADMIN=CAP_NET_ADMIN
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
